# deploy/systemd/wttj-scheduler.timer
[Unit]
Description=Wake WTTJ local scheduler every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=wttj-scheduler.service

[Install]
WantedBy=timers.target
